# configs/launchd/com.barbos.meridian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.barbos.meridian</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source "$HOME/.nvm/nvm.sh" &amp;&amp; exec meridian</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MERIDIAN_HOST</key>
        <string>0.0.0.0</string>
        <key>MERIDIAN_PORT</key>
        <string>3456</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/barbos/Library/Logs/meridian.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/barbos/Library/Logs/meridian.err.log</string>
</dict>
</plist>
